# files/gnatsd.conf
start on net-device-up
exec /opt/gnatsd/gnatsd -a localhost -m 8222
